# ci/build-linux-repositories.path
# Trigger the build-linux-repositories when src file paths change

[Path]
PathChanged=TODO

[Install]
WantedBy=multi-user.target
